# scripts/remnawave-node-go.service
# Systemd service file for remnawave-node-go
# Copy to: /etc/systemd/system/remnawave-node-go.service
# Enable:  sudo systemctl enable remnawave-node-go
# Start:   sudo systemctl start remnawave-node-go
# Status:  sudo systemctl status remnawave-node-go
# Logs:    sudo journalctl -u remnawave-node-go -f

[Unit]
# Human-readable description of the service
Description=Remnawave Node Go - Xray VPN Panel Node

# Start this service after the network is available
# This ensures the HTTP server can bind to network interfaces
After=network.target

# Optional: If using a database or other services, add dependencies here
# After=network.target postgresql.service

[Service]
# Simple type: systemd considers the service started immediately after ExecStart
# This is appropriate for long-running processes that don't fork
Type=simple

# User and group to run the service as (create with: sudo useradd -r -s /bin/false remnawave)
# Using a dedicated user improves security by limiting process privileges
User=remnawave
Group=remnawave

# Working directory for the service
# Configuration files and relative paths will be resolved from here
WorkingDirectory=/opt/remnawave

# Load environment variables from file
# This file should contain: SECRET_KEY=<base64-encoded-json-with-certificates>
# Permissions should be: chmod 600 /etc/remnawave/env
EnvironmentFile=/etc/remnawave/env

# Command to start the service
# --config: Path to JSON configuration file with server settings
ExecStart=/usr/local/bin/remnawave-node-go --config /etc/remnawave/config.json

# Restart policy: always restart the service if it exits
# This ensures the service stays running even after crashes
Restart=always

# Time to wait before restarting (prevents rapid restart loops)
RestartSec=5

# Maximum restart attempts within a time window (optional, uncomment if needed)
# StartLimitBurst=5
# StartLimitIntervalSec=60

# Send stdout to systemd journal (viewable via journalctl)
StandardOutput=journal

# Send stderr to systemd journal (viewable via journalctl)
StandardError=journal

# Add identifier for easier log filtering: journalctl -t remnawave-node-go
SyslogIdentifier=remnawave-node-go

# --- Security Hardening ---

# Prevent the service from gaining new privileges via setuid/setgid
NoNewPrivileges=true

# Make the entire filesystem read-only except for /var, /run, /tmp
# Use ReadWritePaths= if the service needs to write to specific locations
ProtectSystem=strict

# Make /home, /root, and /run/user inaccessible
ProtectHome=true

# Create a private /tmp directory for this service
# Prevents information leakage through shared /tmp
PrivateTmp=true

# Restrict the service from gaining any capabilities
CapabilityBoundingSet=

# Only allow specific capabilities if needed (e.g., for binding to ports < 1024)
# AmbientCapabilities=CAP_NET_BIND_SERVICE

# Restrict address families to only those needed (IPv4, IPv6, Unix sockets)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Prevent loading kernel modules
ProtectKernelModules=true

# Prevent access to kernel tunables (/proc/sys, /sys)
ProtectKernelTunables=true

# Prevent access to kernel logs
ProtectKernelLogs=true

# Protect control groups
ProtectControlGroups=true

# Restrict realtime scheduling (prevents DoS via resource exhaustion)
RestrictRealtime=true

# Prevent memory execution (W^X policy)
MemoryDenyWriteExecute=true

# Lock down personality to native (prevents switching to 32-bit mode)
LockPersonality=true

# Restrict namespace creation
RestrictNamespaces=true

# Private network namespace (uncomment only if service doesn't need network)
# PrivateNetwork=true

# Limit file descriptors (adjust based on expected connections)
LimitNOFILE=65535

# --- Optional Resource Limits ---

# Memory limit (uncomment and adjust as needed)
# MemoryMax=512M

# CPU quota (uncomment and adjust as needed)
# CPUQuota=80%

[Install]
# Start this service when the system reaches multi-user mode (normal boot)
WantedBy=multi-user.target
